# automation/com.curbsiders.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template: run the weekly incremental ingest.

  This is a TEMPLATE. Replace __REPO__ with the absolute path to your checkout
  (e.g. /Users/you/src/curbsiders_to_trials), then install it. See automation/README.md.
  Runs every Sunday at 06:00 local time.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.curbsiders.ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/automation/run_ingest.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO__/ingest.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/ingest.log</string>
</dict>
</plist>
